# polybot.service
[Unit]
Description=PolyBot Polymarket Scanner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/polybot
ExecStart=/home/ubuntu/.local/bin/uv run polybot
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=5

# Environment hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/polybot/logs /home/ubuntu/polybot/data
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polybot

[Install]
WantedBy=multi-user.target
